# supervisor/queue-worker.conf
[program:agencyos-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/agencyos/current/artisan queue:work redis --sleep=3 --max-jobs=100 --tries=3 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/agencyos-queue.log
stdout_logfile_maxbytes=10MB
